# systemd/pubkey-pinning-popup.socket.in
[Socket]
ListenStream=@SOCK_PATH@
SocketMode=0777
Service=pubkey-pinning-popup.service

[Install]
WantedBy=sockets.target
