# systemd/pubkey-pinning-popup.service.in
[Unit]
Description=TPKP(Tizen Public Key Pinning) popup service

[Service]
ExecStart=/usr/bin/tpkp-popup
Type=simple
EnvironmentFile=/run/tizen-system-env
Sockets=pubkey-pinning-popup.socket

[Install]
WantedBy=default.target
