# systemd/pubkey-pinning-popup.socket
[Socket]
ListenStream=/tmp/.tpkp-ui-backend.sock
SocketMode=0777
Service=pubkey-pinning-popup.service

[Install]
WantedBy=sockets.target
